# mongodb_exporter/root/etc/services.d/prom_mongodb_exporter/run
#!/usr/bin/with-contenv sh

exec /bin/mongodb_exporter $EXPORTER_OPTIONS
