# alpine-prom_mysqld_exporter/root/etc/services.d/prom_mysqld_exporter/run
#!/usr/bin/with-contenv sh

exec /bin/mysqld_exporter
