# alpine-dataloop-agent/root/etc/services.d/dataloop-agent/run
#!/usr/bin/with-contenv sh

exec /usr/bin/python /dataloop-agent/agent.py -a $DATALOOP_AGENT_KEY;
