# agent/root/etc/services.d/dataloop-agent/run
#!/bin/sh

if [ ! ${DATALOOP_DEBUG} ] ; then
  DATALOOP_DEBUG="no"
fi

if [ ! ${DATALOOP_SERVER} ] ; then
  DATALOOP_SERVER="wss://agent.dataloop.io"
fi

if [ ! ${DATALOOP_NAME} ] ; then
    DATALOOP_NAME=$(hostname) || 'dataloop'
fi

if [ ! ${DATALOOP_TAGS} ] ; then
    DATALOOP_TAGS='docker'
fi

if [ ${DATALOOP_FINGERPRINT} ]; then
cat <<EOF >>/etc/dataloop/agent.finger
[private]
fingerprint = ${DATALOOP_FINGERPRINT}
EOF
fi

echo "setting debug level"
sed -e 's/<<DEBUG>>/'"${DATALOOP_DEBUG}"'/' -i /etc/dataloop/agent.yaml
echo "sedding dataloop apikey"
sed -e 's/<<APIKEY>>/'"${DATALOOP_AGENT_KEY}"'/' -i /etc/dataloop/agent.yaml
echo "sedding dataloop server"
# use ~ to delimit this sed, as replacement contains http:// slashes
sed -e 's~<<SERVER>>~'"${DATALOOP_SERVER}"'~' -i /etc/dataloop/agent.yaml
echo "sedding dataloop agent protocol"
sed -e 's~http://~ws://~' -i /etc/dataloop/agent.yaml
sed -e 's~https://~wss://~' -i /etc/dataloop/agent.yaml
echo "sedding dataloop name"
sed -e 's/<<NAME>>/'"${DATALOOP_NAME}"'/' -i /etc/dataloop/agent.yaml
echo "sedding dataloop tags"
sed -e 's/<<TAGS>>/'"${DATALOOP_TAGS}"'/' -i /etc/dataloop/agent.yaml

exec /usr/bin/python /opt/dataloop/agent/agent.py --config /etc/dataloop/agent.yaml
